# deploy/systemd/django-6-blog-backup.timer
[Unit]
Description=Daily Django blog backup
[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m
[Install]
WantedBy=timers.target
