# deploy/systemd/django-6-blog-deployment-poller.timer
[Unit]
Description=Schedule GitHub deployment polling for Django 6 Blog

[Timer]
OnBootSec=45s
OnUnitActiveSec=2min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
